# etc/rpi-security.service
[Unit]
Description=rpi-security service
After=multi-user.target

[Service]
Type=idle
StandardOutput=null
TimeoutStartSec=15
ExecStartPre=/sbin/iw phy phy0 interface add mon0 type monitor
ExecStartPre=/sbin/ifconfig mon0 up
ExecStart=/usr/local/bin/rpi-security.py
ExecStop=/usr/bin/pkill rpi-security.py
ExecStopPost=/sbin/iw dev mon0 del

[Install]
WantedBy=multi-user.target
